# app/services/bots/ai-ml-bots/bot-usdjpy-ml-ensemble/bot-usdjpy-ml-ensemble.service
[Unit]
Description=USD/JPY ML Ensemble Auto Trader
After=network.target

[Service]
Type=simple
User=admintrader
WorkingDirectory=/home/admintrader/tradermain/bots/bot-usdjpy-ml-ensemble
Environment="PATH=/home/admintrader/tradermain/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/admintrader/tradermain/venv/bin/python3 /home/admintrader/tradermain/bots/bot-usdjpy-ml-ensemble/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

